# supervisor/firstrun.ini
[program:firstrun]
command=/root/.firstrun/firstrun.sh
priority=4
exitcodes=0

