# supervisor/kibana.ini
[program:kibana]
command=/opt/kibana/bin/kibana.sh
user=nobody
priority=3

