# supervisor/kibana.ini
[program:kibana]
command=/opt/kibana/bin/kibana
user=nobody
priority=3
stdout_logfile=syslog
stderr_logfile=syslog

