# supervisor/rsyslog.ini
[program:rsyslog]
command=/usr/sbin/rsyslogd -n
priority=1

