# supervisor/rsyslog.ini
[program:rsyslog]
command=/usr/sbin/rsyslogd -n
priority=1
stdout_logfile=syslog
stderr_logfile=syslog

