# supervisor/rsyslog.ini
[program:rsyslog]
command=/usr/sbin/rsyslogd -n
priority=2
stdout_logfile=syslog
stderr_logfile=syslog

